# heater_control@.service
[Unit]
Description=heater service control

[Service]
User=pi
Type=idle
Restart=no
Environment="ARGS=%I"
ExecStart=/bin/sh /home/pi/raspb-scripts/main.sh $ARGS
